# .config/systemd/user/udiskie.service
[Unit]
Description=Automounting
After=default.target

[Service]
ExecStart=/usr/bin/udiskie
Restart=always

[Install]
WantedBy=default.target
